# board/jxhmi/common/overlay/etc/profile.d/env.sh
#!/bin/sh

export QT_QPA_PLATFORM=linuxfb
export TSLIB_TSDEVICE=/dev/input/event0
